# deploy/systemd/finny-dashboard.service
[Unit]
Description=Finny Dashboard (TanStack Start server on 127.0.0.1:3001)
After=network-online.target
Wants=network-online.target
# Note: the Hermes gateway runs as a bare process under user `ubuntu`
# (`/home/ubuntu/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run`),
# not a systemd unit. The dashboard talks to it over loopback HTTP.
# If gateway availability becomes brittle, wrap the gateway in its own
# systemd unit and add Requires= here.

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/finny/dashboard
EnvironmentFile=/opt/finny/dashboard/.env
Environment=HOST=127.0.0.1
Environment=PORT=3001
ExecStart=/usr/bin/node /opt/finny/dashboard/server-entry.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
